# crontab
*      *  *   *    *       pkill -USR1 "ocelot"
